# files/etc/profile.d/hadoop.sh
export HADOOP_CONF_DIR=/etc/hadoop
export HADOOP_HOME=/opt/hadoop
export JAVA_HOME=/usr/lib/jvm/java-8-openjdk-amd64/jre
export PATH=$PATH:$HADOOP_HOME/bin:$HADOOP_HOME/sbin
